# packages/network/syncthing/system.d/syncthing.service
[Unit]
Description=Start SyncThing
After=unofficialos.service

[Service]
Type=simple
ExecStart=/usr/bin/start_syncthing.sh

[Install]
WantedBy=unofficialos.service
